# units/pacman-sync/pacman-sync.timer
[Unit]
Description=Sync but do not install Pacman Files

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=multi-user.target
