# units/reboot/reboot.service
[Unit]
Description=Reboot every day at 23:15

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
