# units/shutdown/shutdown.timer
[Unit]
Description=Run shutdown at 8 PM daily

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
